# configs/supervisor/supervisord.conf

[supervisord]
nodaemon=true
user=root

[inet_http_server]
port=127.0.0.1:7000
username = xadmin
password = xadmin


[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[program:xview]
directory=/x-view
command=/usr/local/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 backend.wsgi --log-level info
autostart=true
autorestart=true
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout


[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=5
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout



[program:xray]
command=/usr/local/bin/xray run -c /x-view/config.json
autostart=true
autorestart=true
startretries=5
startsecs=5
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout